# packaging/hubproxy.openrc
#!/sbin/openrc-run

name="hubproxy"
description="Docker and GitHub acceleration proxy server"
command="/usr/bin/hubproxy"
command_user="hubproxy:hubproxy"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/etc/hubproxy/hubproxy.stdout.log"
error_log="/etc/hubproxy/hubproxy.stderr.log"
supervisor="supervise-daemon"
respawn_delay=5
respawn_max=0

export CONFIG_PATH="/etc/hubproxy/config.toml"

depend() {
    need net
    after firewall
}
